# ansible/computer_room/mirror_slackware/files/etc/cron.daily/mirror-slackware-job.sh
#!/bin/env bash
LOGFILE="/var/log/mirror.log"
date >> "$LOGFILE"
/bin/bash /usr/local/bin/mirror-slackware-current.sh &>> "$LOGFILE"

