# ansible/computer_room/packages_server/files/etc/profile.d/spack_setup.sh
#!/bin/bash
if [ -f /packages/spack/share/spack/setup-env.sh ]; then
    source /packages/spack/share/spack/setup-env.sh
fi
